# Funry_Switch_Gateway/rootfs/etc/services.d/funry_switch_gateway/run
#!/usr/bin/with-contenv bashio

exec python /src/main.py
